# supervisord.conf
[supervisord]
nodaemon=true

[program:fhem]
command=/etc/init.d/fhem start

# autostart for homebridge disabled
# [program:homebridge]
# command=/usr/local/bin/homebridge start