# supervisord.conf
[supervisord]
nodaemon=true

[program:fhem]
command=/etc/init.d/fhem start

# autostart for homebridge currently disabled
# [program:homebridge]
# command=/home/homebridge/bin/homebridge start